# Source/Infrastructure/Hosting/Hosts/EtAlii.Ubigia.Infrastructure.Hosting.DockerHost/supervisord.conf
[supervisord]
nodaemon=true

# [program:api]
# command=/usr/bin/dotnet /app/dotnetcore-ssh.dll

[program:ssh]
command=/usr/sbin/sshd -D